# llm-clients/anythingllm/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:fcgiwrap]
command=/usr/sbin/fcgiwrap -s unix:/var/run/fcgiwrap.socket
autostart=true
autorestart=true
redirect_stderr=true
